# roles/ganeti-hoogle/files/etc/cron.hourly/hoogle-update
#!/bin/sh
make -C /srv/hoogle >/srv/hoogle/update.log 2>&1
